# scripts/systemd/kater-cloudflared.service
[Unit]
Description=Cloudflare Tunnel for Kater
After=network-online.target kater.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/kater/tunnel.env
ExecStart=/usr/local/bin/cloudflared tunnel run --token ${KATER_TUNNEL_TOKEN}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
