# scripts/systemd/kater.service.example
[Unit]
Description=Kater MCP Gateway (secured)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/OrgChefgroep/kater-dev-tools
EnvironmentFile=-%h/.config/kater/kater.env
Environment=KATER_PUBLIC=1
Environment=KATER_AUTH_MODE=oauth
Environment=KATER_RATE_LIMIT=60
Environment=KATER_HOST=127.0.0.1
# Set your public hostname (example — replace before deploy):
Environment=KATER_DOMAIN=kater.example.com
Environment=KATER_CORS_ORIGINS=https://kater.example.com
# Validate config + connectivity before serving (systemd already guarantees a
# single instance, so no pkill is needed).
ExecStartPre=%h/OrgChefgroep/kater-dev-tools/.venv/bin/kater doctor
ExecStart=%h/OrgChefgroep/kater-dev-tools/.venv/bin/kater serve --profile core --host 127.0.0.1
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15
KillMode=mixed

# Hardening: drop privileges/capabilities and restrict filesystem/process surface.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
UMask=0077
CapabilityBoundingSet=
AmbientCapabilities=
# Kater keeps SQLite state under .kater/ and a runtime dir.
ReadWritePaths=%h/OrgChefgroep/kater-dev-tools/.kater
RuntimeDirectory=kater
StateDirectory=kater

[Install]
WantedBy=default.target
